# generate-rotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sswanepoel.wallpaperchanger</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/sswanepoel/workspace/getsmarter-adhoc/bin/python3</string>
        <string>/Users/sswanepoel/workspace/getsmarter-adhoc/ai/wallpaper-generator/generate-wallpaper.py</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
